# Slax/debian/modules/02-xorg/rootcopy/lib/systemd/system/xorg.service
[Unit]
Description=X-Window
After=systemd-user-sessions.service

[Service]
ExecStart=su --login -c /usr/bin/startx
