# Slax/debian9/rootcopy/lib/udev/rules.d/90-slax-automount.rules
# we don't care about loop* and ram* devices
KERNEL=="[!lr]*", SUBSYSTEM=="block", RUN+="/sbin/slax-automount %r/%k"
